# supervisor/cowrie.conf
[program:cowrie]
command=/home/cowrie/cowrie/bin/cowrie start
directory=/home/cowrie/cowrie/
user=cowrie
autostart=true
autorestart=unexpected
exitcodes=1
redirect_stderr=true
stdout_logfile=NONE
